# projects/ant-worker-node-metrics-exporter/ant-worker-node-metrics-exporter.service
[Unit]
Description=The typesofants per-host metrics exporter!

[Service]
Type=simple
EnvironmentFile={{INSTALL_DIR}}/.env
Environment=TYPESOFANTS_SECRET_DIR={{INSTALL_DIR}}/secrets
ExecStart={{INSTALL_DIR}}/node_exporter --web.listen-address=":{{PORT}}" --collector.systemd
WorkingDirectory={{INSTALL_DIR}}
Slice=typesofants.slice
Restart=always

[Install]
WantedBy=multi-user.target
